# rootfs/etc/services.d/frenchstream/finish
#!/usr/bin/with-contenv bash
exec /run/s6/basedir/bin/halt
