# rootfs/etc/services.d/frenchstream/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting FrenchStream Stremio addon..."

cd /app
exec node index.js
