# rootfs/etc/services.d/frenchstream/run
#!/usr/bin/with-contenv bash
exec node /app/index.js
